# files/etc/uci-defaults/zzz-fullcone
#!/bin/sh
uci set firewall.@defaults[0].fullcone="1"
uci commit firewall
/etc/init.d/dnsmasq restart
/etc/init.d/firewall restart
exit 0
